# scripts/civil-service-assistant.service
[Unit]
Description=Civil Service Assistant
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/civil-service-assistant
Environment="COZE_WORKSPACE_PATH=/opt/civil-service-assistant"
ExecStart=/usr/bin/python3 /opt/civil-service-assistant/src/main.py
Restart=on-failure
RestartSec=5

# 日志配置
StandardOutput=append:/var/log/civil-service-assistant.out.log
StandardError=append:/var/log/civil-service-assistant.err.log

[Install]
WantedBy=multi-user.target
